# systemd-timer/templates/etc_systemd_system_timer-name.timer.j2
[Unit]
Description={{ systemd_timer_description }}

[Timer]
OnCalendar={{ systemd_timer_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
