# systemd-service/templates/etc_systemd_system_service-name.service.j2
[Unit]
Description={{ systemd_service_description }}

[Service]
Type={{ systemd_service_type }}
ExecStart={{ systemd_service_command }}

[Install]
WantedBy=multi-user.target
